# container_src/config/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chown=root:www
chmod=0770

[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
nodaemon=true
pidfile=/var/run/supervisord.pid
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:django-app-init]
command="/init-django-app.sh"
user=www
priority=1
autostart=true
autorestart=false
stdout_logfile=/var/log/supervisor/%(program_name)s.log
redirect_stderr=true

[program:django-app-run]
command="/run-django-app.sh"
user=www
priority=2
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
redirect_stderr=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
redirect_stderr=true
